# systemd/updatenodelist.timer
[Unit]
Description=MMDVMHost 60 seconds after boot

[Timer]
OnBootSec=60sec

[Install]
WantedBy=multi-user.target
